# Library/Ashmem/ld.so.preload
/usr/local/lib/libandroid-shmem-gnueabihf.so